# launchd/com.aichan.killswitch-drill.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.aichan.killswitch-drill</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>cd /Users/fujihiranoborudai/Downloads/agent/ai-chan && bash scripts/killswitch_drill.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Day</key>
    <integer>1</integer>
    <key>Hour</key>
    <integer>3</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/fujihiranoborudai/Downloads/agent/ai-chan/logs/killswitch_drills/launchd.out</string>
  <key>StandardErrorPath</key>
  <string>/Users/fujihiranoborudai/Downloads/agent/ai-chan/logs/killswitch_drills/launchd.err</string>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
