# launchd/com.aichan.dr-drill.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.aichan.dr-drill</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>cd /Users/fujihiranoborudai/Downloads/agent/ai-chan &amp;&amp; bash scripts/dr_drill.sh</string>
  </array>
  <!-- 四半期 1 回: 1 / 4 / 7 / 10 月 1 日 04:00 JST -->
  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Month</key>
      <integer>1</integer>
      <key>Day</key>
      <integer>1</integer>
      <key>Hour</key>
      <integer>4</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Month</key>
      <integer>4</integer>
      <key>Day</key>
      <integer>1</integer>
      <key>Hour</key>
      <integer>4</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Month</key>
      <integer>7</integer>
      <key>Day</key>
      <integer>1</integer>
      <key>Hour</key>
      <integer>4</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Month</key>
      <integer>10</integer>
      <key>Day</key>
      <integer>1</integer>
      <key>Hour</key>
      <integer>4</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
  </array>
  <key>StandardOutPath</key>
  <string>/Users/fujihiranoborudai/Downloads/agent/ai-chan/logs/dr/launchd.out</string>
  <key>StandardErrorPath</key>
  <string>/Users/fujihiranoborudai/Downloads/agent/ai-chan/logs/dr/launchd.err</string>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
